# snapshot-weekly.service
[Unit]
Description=takes weekly btrfs snapshots of root and home directories
After=default.target
Wants=default.target

[Service]
ExecStart=snapshot-weekly.sh
