# systemd/snapshot-hourly.service
[Unit]
Description=takes hourly btrfs snapshots of root and home directories
After=default.target
Wants=default.target

[Service]
ExecStart=snapshot-hourly.sh
